# itamae/hosts/sakura.nna774.net/files/etc/systemd/system/certbot-renew.timer
[Unit]
Description=Run certbot-renew

[Timer]
OnBootSec=15min
OnUnitActiveSec=10d

[Install]
WantedBy=timers.target
