# itamae/cookbooks/mackerel-agent/templates/etc/systemd/system/mackerel-agent.service
[Unit]
Description=Mackerel Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/mackerel-agent
Restart=always
PrivateTmp=true

[Install]
WantedBy=multi-user.target
