# itamae/hosts/ushio.nna774.net/files/etc/systemd/system/emerge-sync.timer
[Unit]
Description=emerge sync timer

[Timer]
OnBootSec=1d
OnUnitActiveSec=2d

[Install]
WantedBy=timers.target
