# itamae/hosts/ushio.nna774.net/files/etc/systemd/system/emerge-sync.service
[Unit]
Description=emerge --sync

[Service]
ExecStart=/usr/bin/eix-sync
Type=oneshot

[Install]
WantedBy=multi-user.target
