# itamae/hosts/ushio.compute.kitashirakawa.dark-kuins.net/files/etc/systemd/system/emerge-sync.service
[Unit]
Description=emerge --sync
After=network.target

[Service]
ExecStart=/usr/bin/eix-sync
Type=oneshot

[Install]
WantedBy=multi-user.target
