# itamae/hosts/ushio.compute.kitashirakawa.dark-kuins.net/files/etc/systemd/system/update-ddns.service
[Unit]
Description=update dns

[Service]
ExecStart=/usr/local/bin/update-ddns
Type=oneshot
User=nobody

[Install]
WantedBy=multi-user.target
